# scripts/wine-mobile-sync.timer
[Unit]
Description=Weeno mobile IPA sync ~2 min (comme beer-mobile-sync)

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
